# packaging/templates/pitrac.service
[Unit]
Description=PiTrac Launch Monitor
After=network.target activemq.service tomee.service
Wants=activemq.service tomee.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi
Environment="HOME=/home/pi"
Environment="USER=pi"
Environment="LD_LIBRARY_PATH=/usr/lib/pitrac"
Environment="PITRAC_BASE_IMAGE_LOGGING_DIR=/home/pi/LM_Shares/Images/"
Environment="PITRAC_WEBSERVER_SHARE_DIR=/home/pi/LM_Shares/WebShare/"
Environment="PITRAC_MSG_BROKER_FULL_ADDRESS=tcp://localhost:61616"
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/pitrac run --foreground --system_mode=camera1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target